# infra/deploy/teatiers-backup.timer
[Unit]
Description=Run the TeaTiers Postgres backup daily

[Timer]
# 03:17 UTC daily, with catch-up if the VM was off at the scheduled time.
OnCalendar=*-*-* 03:17:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
