# infra/deploy/teatiers-backup.service
[Unit]
Description=TeaTiers Postgres logical backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# backup.sh reads /opt/teatiers/.env for DB name/user. Off-box upload is enabled by dropping a
# /opt/teatiers/backup.env with BACKUP_S3_URI + AWS_ACCESS_KEY_ID/AWS_SECRET_ACCESS_KEY (from the
# `tofu output backup_*`); the leading `-` keeps the unit working when that file is absent (#77).
EnvironmentFile=-/opt/teatiers/backup.env
ExecStart=/opt/teatiers/backup.sh
